# modules/trackmania/files/rc.local
#!/bin/sh -e
#
# rc.local
#
(
cd /opt/TmDedicatedServer
/opt/TmDedicatedServer/TrackmaniaServer /dedicated_cfg=dedicated_cfg.txt /lan /game_settings=MatchSettings/Nations/NationsBlue.txt
)
exit 0
.
