# modules/trackmania/files/rc.local
#!/bin/sh -e
#
# rc.local
#
(
cd /opt/TmDedicatedServer
/opt/TmDedicatedServer/TrackManiaServer /game=nations /lan game_settings=MatchSettings/Lan/AdvancedTraining.txt /dedicated_cfg=dedicated.cfg
)
exit 0
