# modules/factorio/files/rc.local
#!/bin/sh -e
#
# rc.local
#
tmux new-session -d "/opt/factorio/bin/x64/factorio"
exit 0
