# modules/openttd/files/rc.local
#!/bin/sh -e
#
# rc.local
#
tmux new-session -d openttd -D
exit 0
