# modules/openttd/files/rc.local
#!/bin/sh -e
#
# rc.local
#
tmux new-session -d "/opt/openttd/openttd -D"
exit 0
