# modules/terraria/files/rc.local
#!/bin/sh -e
#
# rc.local
#
tmux new-session -d "mono /opt/terraria/TerrariaServer.exe -autocreate 1 -world BolkLan"
exit 0
