# modules/terraria/files/rc.local
#!/bin/sh -e
#
# rc.local
#
tmux new-session -d "/opt/terraria/tshock/TerrariaServer.exe"
exit 0
